# src/com.claude.sync-theme.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.sync-theme</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/laura/.claude/sync-theme.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/laura/Library/Preferences/.GlobalPreferences.plist</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-sync-theme.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-sync-theme.log</string>
</dict>
</plist>